# rva.service
[Unit]
Description=RVA Website

[Service]
User=rva
Group=rva
WorkingDirectory=/home/rva/RVA/current

ExecStart=/bin/bash -lc 'docker-compose up --build -d'
